# scripts/v1.5/ota-client.service
[Unit]
Description=OTA Client Daemon
After=network.target

# 配置见 -config 指向的 YAML（config_url、version_file、agent_id、processes 等）。

[Service]
Type=simple
User=root
ExecStart=/home/arenatech/client/bin/ota-agent -config=/home/arenatech/client/bin/ota-agent.yaml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
