# scripts/v2.0/ota-client.service
[Unit]
Description=OTA Client Daemon
After=network.target

# 可选远程日志：在 YAML 中设置 log_upload.enabled: true 等，见 proj-ota/LOG_UPLOAD.md

[Service]
Type=simple
User=root
ExecStart=/home/arenatech/client/bin/ota-agent -config=/home/arenatech/client/bin/ota-agent.yaml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
