# ota-agent/agent.service
[Unit]
Description=OTA Agent Daemon
After=network.target

# 在 -config 指向的 YAML 中填写 config_url、version_file、agent_id、processes 等；
# 相对路径 version_file 相对于 ota-agent 二进制所在目录。

[Service]
Type=simple
User=root
ExecStart=/home/arenatech/client/bin/ota-agent -config=/home/arenatech/client/bin/ota-agent.yaml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
